# cameraapp/cameraapp.service
[Unit]
Description=chs camera app

[Service]
Type=simple
ExecStart=/usr/bin/cameraapp

[Install]
WantedBy=multi-user.target
